# overlays/growfs/files/growfs.service
[Unit]
Description=Grow root filesystem on first boot
DefaultDependencies=no
After=local-fs.target
RequiresMountsFor=/
ConditionPathExists=!/var/lib/tachyon/growfs.done

[Service]
Type=oneshot
ExecStart=/usr/bin/growfs.sh
RemainAfterExit=yes

[Install]
# If you adopt this, you can remove the two ln -sfn steps from the overlay JSON
WantedBy=multi-user.target ffbm.target